# crontab
0 18 * * * /usr/local/bin/ffmpeg -f avfoundation -framerate 30 -i 0 -frames:v 1 -y ~/files/img/6pm/"$(date +\%s)".jpg
